# files/etc/uci-defaults/20_firewall_disable
#!/bin/sh

exit 0 

service firewall stop
service firewall disable
exit 0
