# files/etc/uci-defaults/20_usbnet_enable
#!/bin/sh

[ "`uci get network.@device[0].ports`" == "usb0" ] && exit 0 # only enable if disabled

uci set network.@device[0].ports='usb0'

uci commit network

exit 0
