# files/etc/rc.local
# This file is from the default image to store audio state

modprobe g_ether

exit 0
