# files/etc/rc.local
# This file is from the default image

[ -e /tmp/reboot ] && reboot & 

modprobe g_ether

exit 0
